# bundle/docker/etc/init.d/sidekiq
#!/sbin/openrc-run

depend() {
    need net
}

command="bundle exec sidekiq"
command_background=1
command_user="decidim"
command_directory="/home/decidim/app"
command_logfile="/home/decidim/app/log/run.sidekiq.log"
command_logfile_maxbytes=4MB
command_logfile_backups=10
command_stderr="/home/decidim/app/log/errors.sidekiq.log"
